# files/nagios3.cron
#!/bin/sh -e
find /var/log/nagios3/archives -ctime +90 | xargs -r rm
